# docker/supervisord-secure-sync.conf
[program:secure-sync]
command=bin/secure-sync.sh
directory=/home/app/webapp
autorestart=false
startsecs=0
priority=1
stdout_logfile=syslog
stderr_logfile=syslog